# amlogic-s9xxx/platform-files/amlogic/rootfs/etc/inittab
platform-files/amlogic/rootfs/etc/inittab
